# service/crow-server/run
#!/bin/sh

dir=$(pwd)
dir=$(dirname "$dir")
dir=$(dirname "$dir")
exec 2>&1
exec "$dir/server/crow-server"
